# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm
user=www-data

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=www-data
